# systemd/kahu-daemon.service
[Unit]
Description=KAHU vessel radar daemon
Documentation=https://github.com/KAHU-radar/kahu-vessel-rs
After=network-online.target mayara-server.service
Wants=network-online.target
Requires=mayara-server.service

[Service]
EnvironmentFile=/etc/default/kahu
ExecStart=/usr/local/bin/kahu-daemon \
    --ws-url ws://localhost:6502/signalk/v2/api/vessels/self/radars/${RADAR_ID}/spokes \
    --api-key ${KAHU_API_KEY} \
    --land-filter \
    --spoke-timeout 60
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
